# .config/profile.d/env
WGET2RC="${XDG_CONFIG_HOME}/wget2/config"
export WGET2RC

WINEPREFIX="${XDG_DATA_HOME}/wine"
export WINEPREFIX

CMAKE_LLVM="-DCMAKE_C_COMPILER=clang -DCMAKE_CXX_COMPILER=clang++ -DLLVM_ENABLE_LLD=ON -DLLVM_ENABLE_IDE=ON"
export CMAKE_LLVM

CPP_BASE="-O3 -std=c++2a"
CPP_STRICT_ARGS="-Wall -Warray-bounds=2 -Wdouble-promotion -Weffc++ -Wextra -Wimplicit-fallthrough=5 -Winit-self -Wfloat-equal -Wformat=2 -Wformat-overflow=2 -Wformat-signedness -Wformat-truncation=2 -Wmisleading-indentation -Wmissing-include-dirs -Wnull-dereference -Wparentheses -Wpedantic -Wshadow -Wstrict-aliasing=1 -Wstrict-overflow=5 -Wstringop-overflow=4 -Wtautological-compare -Wuninitialized -Wvector-operation-performance"
CPP_NO_RTTX="-fno-rtti -fno-exceptions"
CPP_INCLUDE="-I include -Wno-missing-include-dirs"
CPP_ARGS="${CPP_BASE} ${CPP_NO_RTTX} ${CPP_STRICT_ARGS} ${CPP_INCLUDE}"
CPP_A_ARGS="${CPP_BASE} ${CPP_NO_RTTX} ${CPP_INCLUDE} -S -fno-ident -finhibit-size-directive -fno-PIC -fno-asynchronous-unwind-tables -fno-stack-protector -fcf-protection=none -fno-dwarf2-cfi-asm -o-"
export CPP_BASE
export CPP_STRICT_ARGS
export CPP_NO_RTTX
export CPP_INCLUDE
export CPP_ARGS
export CPP_A_ARGS

DPREFIX="${XDG_DATA_HOME}/darling"
export DPREFIX

ANDROID_HOME="${XDG_LOCAL_HOME}/lib/android"
ANDROID_SDK_HOME="${ANDROID_HOME}"
ANDROID_NDK_ROOT="${ANDROID_HOME}/ndk-bundle"
export ANDROID_HOME
export ANDROID_SDK_HOME
export ANDROID_NDK_ROOT

XCODE_APP="/home/maze/Desktop/Xcode-beta.app"
export XCODE_APP

LOCINCLUDE="-fPIC -I${XDG_LOCAL_HOME}/lib/Qt/5.14.1/gcc_64/include -I${XDG_LOCAL_HOME}/lib/Qt/5.14.1/gcc_64/include/QtCore -I${XDG_LOCAL_HOME}/lib/Qt/5.14.1/gcc_64/include/QtBluetooth"
LOCFILE="${LOCINCLUDE}"
export LOCINCLUDE
export LOCFILE
